# ndr-tools/services/ndr-api.service
[Unit]
Description=NDR API (Django)
Documentation=https://github.com/your-org/ndr
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/tools/pcapql/ndr-api/flowtracex_api
Environment=APP_MODE=production
Environment=DJANGO_SETTINGS_MODULE=flowtracex_api.settings
EnvironmentFile=-/opt/tools/pcapql/.env
ExecStart=/usr/bin/python3 manage.py runserver 0.0.0.0:8000
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ndr-api

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
